# solar-miner.service
[Unit]
Description=Solar Miner Controller
After=network.target modbus-sniffer.service
Wants=modbus-sniffer.service

[Service]
Type=simple
User=jstormes
WorkingDirectory=/home/jstormes/solar1
ExecStart=/usr/local/bin/node index.js
Restart=always
RestartSec=10

# Site-specific config lives in /etc/solar-miner.env (not in repo)
EnvironmentFile=/etc/solar-miner.env

[Install]
WantedBy=multi-user.target
